# antigravity-bridge.service
[Unit]
Description=Antigravity Signal Mailbox Bridge
After=signal-daemon.service

[Service]
Type=simple
User=<user>
WorkingDirectory=/<user>/antigravity_mailbox
EnvironmentFile=/<user>/antigravity_mailbox/.env
ExecStart=/<user>/antigravity_mailbox/venv/bin/python /<user>/antigravity_mailbox/signal_bridge.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
